# deploy/showrunner.service
[Unit]
Description=Showrunner TTRPG engine
After=network.target

[Service]
Type=simple
User=showrunner
WorkingDirectory=/opt/showrunner
ExecStart=/opt/showrunner/venv/bin/uvicorn showrunner.web.app:create_server_app \
    --factory \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 1
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
